# user/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:Xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1024x768x24
autostart=true
autorestart=true

[program:tightvncserver]
command=/usr/bin/tightvncserver :1 -geometry 1024x768 -depth 24 -dpi 96
autostart=true
autorestart=true

